# i3/systemd/twmnd.service
Description=Launch twmn notification daemon

[Service]
ExecStart=/usr/bin/twmnd

[Install]
WantedBy=default.target
